# packages/ui/rocknix-emulationstation/system.d/es.service
[Unit]
Description=EmulationStation-service-for-any-compositor
PartOf=graphical.target

[Service]
Environment=HOME=/storage
Type=simple
ExecStart=/usr/bin/start_es.sh
WorkingDirectory=/storage
Restart=always
RestartSec=2

[Install]
WantedBy=rocknix.target
