# packages/network/syncthing/system.d/syncthing.service
[Unit]
Description=Start SyncThing
After=rocknix.service

[Service]
Type=simple
ExecStart=/usr/bin/start_syncthing.sh

[Install]
WantedBy=rocknix.service
