# packages/wayland/util/waybar/system.d/waybar.service
[Unit]
Description=Waybar - Highly customizable Wayland bar for Sway and Wlroots-based compositors
Requires=sway.service
After=sway.service

[Service]
ExecStart=/usr/bin/waybar
Restart=always
KillSignal=SIGQUIT
RestartSec=3
StartLimitInterval=0

[Install]
WantedBy=rocknix.target