# projects/Qualcomm/devices/SD865/filesystem/usr/lib/udev/rules.d/10-set-mac.rules
ACTION=="add", SUBSYSTEM=="net", ENV{INTERFACE}=="wlan0", RUN+="/usr/bin/set-mac"
