# projects/Qualcomm/devices/SD865/filesystem/usr/lib/udev/rules.d/99-set-mac.rules
ACTION=="add", SUBSYSTEM=="net", ENV{INTERFACE}=="wlan0", RUN+="/usr/bin/set-mac wifi"
ACTION=="add", SUBSYSTEM=="bluetooth", KERNEL=="hci0", RUN+="/storage/set-mac bluetooth"
