# packages/rocknix/udev.d/80-usbgadget.rules
SUBSYSTEMS=="extcon", ACTION=="change", ENV{STATE}=="USB=1", ATTR{cable.0/name}=="USB", ATTR{cable.0/state}=="1", RUN+="/usr/bin/usbgadget trigger"
SUBSYSTEM=="net", ACTION=="add", ENV{DEVTYPE}=="gadget", RUN+="/usr/sbin/ip link set %k up", RUN+="/usr/sbin/brctl addif usbnet %k"
