# packages/rocknix/udev.d/80-usbgadget.rules
# trigger for activating gadget after cable is plugged
SUBSYSTEMS=="extcon", ACTION=="change", RUN+="/usr/bin/usbgadget trigger"
