# packages/wayland/compositor/sway/system.d/sway-touch.service
[Unit]
Description=Set up touchscreen devices on sway
After=essway.service

[Service]
Type=simple
ExecStart=/usr/bin/sway-touch.sh
TimeoutStartSec=0

[Install]
WantedBy=rocknix.target
