# packages/rocknix/system.d/rocknix-report-stats.timer
[Unit]
Description=Timer to report anonymous ROCKNIX install statistics

[Timer]
OnBootSec=5min

[Install]
WantedBy=timers.target
