# packages/rocknix/system.d/rocknix-report-stats.service
[Unit]
Description=Report anonymous ROCKNIX install statistics

[Service]
Type=oneshot
ExecStart=/usr/bin/rocknix-report-stats

[Install]
WantedBy=multi-user.target
