# packages/kernel/linux/udev.d/99-joystick-calibration.rules
ACTION=="add", KERNEL=="input3"

ATTRS{phys}=="retrogame_joypad/input0", RUN+="evdev-joystick -e /dev/input/by-path/platform-rocknix-singleadc-joypad-event-joystick -m -1800 -M 2400 -a 0 -d 16 -f 32; evdev-joystick -e /dev/input/by-path/platform-rocknix-singleadc-joypad-event-joystick -m -1800 -M 1800 -a 1 -d 16 -f 32; evdev-joystick -e /dev/input/by-path/platform-rocknix-singleadc-joypad-event-joystick -m -1800 -M 1800 -a 2 -d 16 -f 32; evdev-joystick -e /dev/input/by-path/platform-rocknix-singleadc-joypad-event-joystick -m -1800 -M 1800 -a 3 -d 16 -f 32"
